# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:api]
directory=/app
command=/bin/sh -c "python scripts/wait_for_services.py && uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 2"
autorestart=true
startretries=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:frontend]
directory=/app/frontend
command=/usr/local/bin/npm start -- -p 3000 -H 0.0.0.0
environment=OLLAMA_GATEWAY_URL="http://127.0.0.1:8000/api/v1",NEXT_PUBLIC_API_BASE_URL="/api/gateway"
autorestart=true
startretries=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
